# packages/agent-soul-framework/connectors/feishu/systemd/hunqi-core@.service
[Unit]
Description=魂器核心 (Hunqi Core) - OpenCode + Soul
# 注意：install-systemd.sh 安装时会将 %h 替换为实际家目录。
# 系统级 systemd 中 %h 展开为 /root（而非 User= 指定的用户），不可直接使用。
Documentation=https://github.com/NeoMei/agent-soul-framework
After=network-online.target
Wants=network-online.target

# 挂起后重启
After=suspend.target hibernate.target hybrid-sleep.target

[Service]
Type=simple
User=%I
Group=%I

# 工作目录
WorkingDirectory=%h/.hunqi/agent-soul-framework/packages/agent-soul-framework

# 加载环境变量
Environment="HOME=%h"
Environment="USER=%I"

# 从 .env 文件加载环境变量
EnvironmentFile=-%h/.hunqi/.env

# 启动命令 - 使用 wrapper 解析 opencode 位置并启动 serve
ExecStart=/bin/sh -c 'OPENCODE=$(%h/.hunqi/bin/opencode-feishu --resolve-opencode 2>/dev/null || command -v opencode); exec "$OPENCODE" serve --port 19876'

# 停止命令（精确匹配端口，避免误杀）
ExecStop=/bin/sh -c 'PID=$(ss -tlnp 2>/dev/null | grep ":19876 " | sed -n "s/.*pid=\\([0-9]*\\).*/\\1/p"); [ -n "$PID" ] && kill "$PID" || true'

# 重启策略
Restart=always
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hunqi-core

# 进程管理
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
